# overthewall/ipsec/supervisord.conf
[supervisord]
nodaemon=true

[program:ipsec]
command=/usr/sbin/ipsec start --nofork

[program:xl2tpd]
command=/usr/sbin/xl2tpd -D

